# supervisor.conf
[supervisord]
nodaemon=true

# 启动Redis
[program:redis]
command=redis-server
autostart=true
autorestart=true
stdout_logfile=/var/log/redis.log
stderr_logfile=/var/log/redis.err.log

# 启动MySQL
[program:mysql]
command=/usr/sbin/mysqld --user=mysql
autostart=true
autorestart=true
stdout_logfile=/var/log/mysql.log
stderr_logfile=/var/log/mysql.err.log
user=root
priority=1

# MySQL初始化（只执行一次）
[program:mysql_init]
command=/bin/bash -c "sleep 5 && if [ ! -f /var/lib/mysql/.initialized ]; then mysql -uroot < /app/mysql-init.sql && mysql -uroot -p${MYSQL_ROOT_PASSWORD} < /app/init_database.sql && touch /var/lib/mysql/.initialized; fi"
autostart=true
autorestart=false
startsecs=0
startretries=0
stdout_logfile=/var/log/mysql_init.log
stderr_logfile=/var/log/mysql_init.err.log
user=root
priority=2

# 启动Spring Boot后端
[program:osprey_backend]
command=/bin/bash -c "sleep 15 && java -jar /app/osprey_be/osprey_be.jar"
autostart=true
autorestart=true
stdout_logfile=/var/log/backend.log
stderr_logfile=/var/log/backend.err.log
priority=3

# 启动Nginx
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
user=root
stdout_logfile=/var/log/nginx.log
stderr_logfile=/var/log/nginx.err.log
priority=4
